# cron/crontab
# m h  dom mon dow   command
*/5 * * * * echo `date`>>/tmp/cronz.log

